# packaging/debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|purge)
        # Stop the daemon for all logged-in users
        for user_dir in /run/user/*; do
            uid=$(basename "$user_dir")
            sudo -u "#$uid" systemctl --user stop pilot-daemon 2>/dev/null || true
            sudo -u "#$uid" systemctl --user disable pilot-daemon 2>/dev/null || true
        done

        # Remove systemd user service
        rm -f /usr/lib/systemd/user/pilot-daemon.service
        systemctl --global daemon-reload 2>/dev/null || true

        # Remove symlink
        rm -f /usr/local/bin/pilot-daemon

        # Remove PolicyKit policy
        rm -f /usr/share/polkit-1/actions/com.pilot.policy

        # Remove desktop file
        rm -f /usr/share/applications/pilot.desktop
        update-desktop-database /usr/share/applications 2>/dev/null || true

        # Remove venv and installation
        rm -rf /opt/pilot

        echo "Pilot removed."
        echo "User data in ~/.config/pilot and ~/.local/share/pilot was preserved."
        echo "To remove user data: rm -rf ~/.config/pilot ~/.local/share/pilot ~/.local/state/pilot"
        ;;
esac

exit 0
